# crontabs
0 5 * * * /usr/local/bin/python /usr/src/app/manage.py extract_products >> /var/log/cron.log 2>&1
